# board/raspberrypi/fs_server/etc/systemd/system/otp-tunnel-server.service
[Unit]
Description=udptunnel server
After=multi-user.target

[Service]
WorkingDirectory=/opt/tunnel/
KillMode=process
ExecStart=server -i server.ini
Restart=always
RestartSec=5
TimeoutStartSec=5

